# .claude/agent-team/supervisor/launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zilliz.doc-agent-approval-consumer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/ABSOLUTE/PATH/feishu-markdown-bridge/.claude/agent-team/bin/doc-agent-approval-consumer.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/feishu-markdown-bridge</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GITHUB_TOKEN</key>
    <string>REPLACE_WITH_TOKEN_OR_USE_LAUNCHD_ENV</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
